# initramfs/initramfs.org/lib/systemd/system/whoopsie.path
[Unit]
Description=Start whoopsie on modification of the /var/crash directory

[Path]
PathChanged=/var/crash

[Install]
WantedBy=multi-user.target
