# initramfs/initramfs-6.14.8/lib/systemd/user/plasma-kwin_wayland.service
[Unit]
Description=KDE Window Manager
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/kwin_wayland_wrapper --xwayland
BusName=org.kde.KWinWrapper
Slice=session.slice
