# initramfs/initramfs-6.14.8/lib/systemd/user/xfce4-session.service
[Unit]
Description=Xfce Session Manager
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/startxfce4
