# initramfs/initramfs-6.14.8/lib/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver
ConditionUser=!@system

[Service]
ExecStart=xscreensaver

[Install]
WantedBy=default.target
