# initramfs/initramfs-6.14.8/lib/systemd/system/uwsgi-app@.socket
[Unit]
Description=Socket for uWSGI app %i

[Socket]
ListenStream=/var/run/uwsgi/%i.socket
SocketUser=www-data
SocketMode=0600

[Install]
WantedBy=sockets.target
